# package/webinterface-localhost/webinterface-localhost.socket
[Socket]
Description=View the web interface if running, on localhost
ListenStream=127.0.0.1:80
FreeBind=true

[Install]
WantedBy=sockets.target
