# package/remux/prerm
#!/usr/bin/env bash

echo "Disabling remux"
systemctl disable --now remux
