# package/draft/prerm
#!/usr/bin/env bash

echo "Disabling draft and restoring xochitl"
systemctl disable draft
systemctl stop draft
systemctl enable xochitl
systemctl start xochitl
