# package/draft/prerm
#!/usr/bin/env bash

echo "Disabling draft"
systemctl disable --now draft
