# package/oxide/postinst
#!/usr/bin/env bash

echo "Disabling xochitl and enabling oxide"
systemctl daemon-reload
systemctl disable xochitl
systemctl stop xochitl
systemctl enable oxide
systemctl start oxide
