# package/oxide/prerm
#!/usr/bin/env bash

echo "Disabling oxide"
systemctl disable --now oxide
