# package/oxide/prerm
#!/usr/bin/env bash

echo "Disabling oxide and restoring xochitl"
systemctl disable oxide
systemctl stop oxide
systemctl enable xochitl
systemctl start xochitl
